# projects/xzenia/recovery/com.xzenia.tier1.recovery.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.xzenia.tier1.recovery</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/marcuscoarchitect/.openclaw/workspace/projects/xzenia/recovery/run_tier1_recovery.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>/Users/marcuscoarchitect/.openclaw/workspace/projects/xzenia/recovery/logs/launchd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/marcuscoarchitect/.openclaw/workspace/projects/xzenia/recovery/logs/launchd.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/marcuscoarchitect/.openclaw/workspace</string>
  </dict>
</plist>
