# projects/xzenia/protocol-node/com.xzenia.protocol-multiplexer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.xzenia.protocol-multiplexer</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/python3</string>
    <string>-m</string>
    <string>uvicorn</string>
    <string>multiplexer_app:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8788</string>
    <string>--no-access-log</string>
    <string>--log-level</string>
    <string>warning</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/marcuscoarchitect/.openclaw/workspace/projects/xzenia/protocol-node</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>/Users/marcuscoarchitect/.openclaw/workspace/projects/xzenia/protocol-node</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/marcuscoarchitect/.openclaw/workspace/projects/xzenia/protocol-node/multiplexer.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/marcuscoarchitect/.openclaw/workspace/projects/xzenia/protocol-node/multiplexer.err.log</string>
</dict>
</plist>
